# doc/service/nextjs-hn-pb.service
[Unit]
Description=Next.js Hacker News + PocketBase Web App
After=network.target

[Service]
Type=exec
User=www-data
WorkingDirectory=/var/www/HN-PB-Next/web
Environment=NODE_ENV=production
ExecStart=/usr/bin/pnpm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nextjs-hn-pb

[Install]
WantedBy=multi-user.target